# build/management/supervisord.conf
[supervisord]
logfile = /tmp/supervisord.log
[program:media-service]
command=npm run start
autostart=true
autorestart=true
stderr_logfile=/var/log/media-service.err.log
stdout_logfile=/var/log/media-service.out.log
# environment=RTSP_STREAM_LIST_ENDPOINT=true
# [program:media-service-ui]
# command=http-server -p 80 /ui/
# autostart=true
# autorestart=true
# stderr_logfile=/var/log/media-service-ui.err.log
# stdout_logfile=/var/log/media-service-ui.out.log
# environment=API_URL=http://127.0.0.1:8080
